# provisioning/salt/roots/salt/monitoring/client/config/bucky.ini
[program:bucky]
command          = /usr/bin/bucky /etc/bucky/bucky.conf
autostart        = true
autorestart      = true
user             = bucky
log_stdout       = true
log_stderr       = true
logfile          = /var/log/bucky/bucky.log
environment  = HOME=/,USER=bucky
